# lib/capistrano/files/puma.service.erb
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
# Foreground process
Type=simple

# Preferably configure a non-privileged user
User=<%= fetch(:deploy_user) %>

# Specify the path to your puma application root
WorkingDirectory=<%= fetch(:deploy_to) %>/current

# Start command
ExecStart=/home/<%= fetch(:deploy_user) %>/.rvm/bin/rvm default do bundle exec puma -C config/puma.rb

Restart=always

[Install]
WantedBy=multi-user.target
